# deploy/nft-okboy-agent-upgrade.timer
[Unit]
Description=Periodic nft-okboy agent self-upgrade

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
